# setup/ansible/roles/install/templates/systemd.service.j2
[Unit]
Description={{service}}
Requires=network-online.target
After=network-online.target

[Service]
{{ serviceLines | default('') }}
WorkingDirectory={{chdir}}
Restart=on-failure
User={{user}}
Group={{user}}
PermissionsStartOnly=true
ExecStart={{execline}}
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target

