# packages/deployment/bigdipper/ag-lcd.service
[Unit]
Description=Agoric Light Client Daemon
Requires=network-online.target
After=network-online.target

[Service]
WorkingDirectory=/home/bigdipper
Restart=on-failure
User=bigdipper
Group=bigdipper
PermissionsStartOnly=true
ExecStart=/home/bigdipper/ag-lcd.sh
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
